# apps/flask-uwsgi-nginx/supervisord/supervisord.conf
[supervisord]

[supervisorctl]

[program:app-uwsgi]
command = /usr/sbin/uwsgi --ini /tmp/uwsgi/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx